# pwsh/Documents/PowerShell/profile.d/atuin.ps1
if (Get-Command atuin -ErrorAction SilentlyContinue) {
    Invoke-Expression (& { (atuin init powershell | Out-String) })
}
