# pwsh/Documents/PowerShell/profile.d/dot.ps1
if ($env:DOTFILES -and (Test-Path "$env:DOTFILES\bin\dot.ps1")) {
    function dot { & "$env:DOTFILES\bin\dot.ps1" @args }
}
