# pwsh/Documents/PowerShell/profile.d/zoxide.ps1
if (Get-Command zoxide -ErrorAction SilentlyContinue) {
    Invoke-Expression (& { (zoxide init powershell --cmd cd | Out-String) })
}
